# raspberry-pi/config/systemd/cupcake-setup.service
[Unit]
Description=CUPCAKE First Boot Setup
After=network.target postgresql.service redis.service
Requires=postgresql.service redis.service
Before=cupcake-web.service

[Service]
Type=oneshot
User=cupcake
Group=cupcake
WorkingDirectory=/opt/cupcake/app
Environment=PATH=/opt/cupcake/venv/bin
Environment=DJANGO_SETTINGS_MODULE=cupcake.settings
Environment=DATABASE_URL=postgresql://cupcake:cupcake123@localhost/cupcake_db
Environment=REDIS_URL=redis://localhost:6379/0
Environment=PYTHONPATH=/opt/cupcake/app
ExecStart=/opt/cupcake/scripts/first-boot-setup.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cupcake-setup

[Install]
WantedBy=multi-user.target
